# user.service
[Unit]
Description = User Service
[Service]
User=roboshop
Environment=MONGO=true
// highlight-start
Environment=REDIS_URL='redis://redis.daws-86.shop:6379'
Environment=MONGO_URL="mongodb://mongodb.daws-86.shop:27017/users"
// highlight-end
ExecStart=/usr/bin/node /app/server.js
SyslogIdentifier=user

[Install]
WantedBy=multi-user.target